# agent/dist/deb/siem-agent_1.0.0_amd64/DEBIAN/postinst
#!/bin/sh
set -e
case "$1" in
  configure)
    # Create service user if missing
    if ! id siem-agent >/dev/null 2>&1; then
      useradd --system --no-create-home --shell /sbin/nologin \
              --home /var/lib/siem-agent siem-agent
    fi
    chown -R siem-agent:siem-agent /var/lib/siem-agent
    chmod 750 /var/lib/siem-agent
    # Config must be writable by agent (stores token after enrollment)
    chown siem-agent:siem-agent /etc/siem-agent/config.yaml
    chmod 600 /etc/siem-agent/config.yaml

    systemctl daemon-reload
    systemctl enable siem-agent.service
    echo "siem-agent installed. Edit /etc/siem-agent/config.yaml then:"
    echo "  sudo systemctl start siem-agent"
    ;;
esac
